# crontab
0 0 * * * root /renew.sh >> /var/log/cron.log 2>&1