# crontab
0 0 * * * /renew.sh >> /var/log/cron.log 2>&1